# packages/pymidi-controller/pymidi_controller-0.3.0.tar.gz/pymidi_controller-0.3.0/pymidi_controller/data/systemd/pymidi-controller.service
[Unit]
Description=🎹 MIDI Controller Background Service
After=graphical-session.target sound.target
Wants=graphical-session.target sound.target

[Service]
Type=simple
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/env pymidi-controller run --mode blocking
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target