# packages/PieData/piedata-0.0.10-py3-none-any.whl/PieData/piedataserver.service
[Unit]
Description=PieData Server (WebSocket API)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env piedataserver
Restart=on-failure
WorkingDirectory=%h
Environment=HOME=%h

[Install]
WantedBy=default.target
